# docker/memcache.ini
[program:memcache]
command=/usr/bin/memcached -p 11211 -u memcached -m 64 -c 1024 -l 127.0.0.1
redirect_stderr=true
stdout_logfile=/home/user/memcache-logs
autorestart=true
